# tradehub/deploy/systemd/teleposter-checker.timer
[Unit]
Description=Почасовая проверка актуальности объявлений teleposter

[Timer]
# Первый запуск через 5 минут после загрузки, далее каждый час
OnBootSec=5min
OnUnitActiveSec=1h
AccuracySec=30s
Persistent=true
Unit=teleposter-checker.service

[Install]
WantedBy=timers.target
