# tradehub/deploy/systemd/teleposter-parser.service
[Unit]
Description=teleposter — синхронизация объявлений из Telegram (parser:sync)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/opt/teleposter/.env
WorkingDirectory=/opt/teleposter
# Пока идёт предыдущий прогон — тихо выходим 0 (без «failed» в journal каждую минуту)
ExecStart=/bin/bash -lc 'exec 200>/run/teleposter-parser.sync.lock; flock -n 200 || exit 0; exec npm run parser:sync'
User=root
Group=root
StandardOutput=journal
StandardError=journal
# Две группы с длинным проходом «вверх» (800ms × много ID) + refresh-images легко >1ч — иначе systemd режет SIGTERM посреди прогона.
TimeoutStartSec=14400

[Install]
WantedBy=multi-user.target
