# tradehub/deploy/systemd/tradehub-parser.service
[Unit]
Description=TradeHub — синхронизация объявлений из Telegram (parser:sync)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/opt/tradehub/.env
WorkingDirectory=/opt/tradehub
# npm из PATH (Node устанавливается вместе с приложением)
ExecStart=/usr/bin/env npm run parser:sync
User=root
Group=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
