# tradehub/deploy/systemd/tradehub.service
[Unit]
Description=TradeHub (SvelteKit / Node)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tradehub
EnvironmentFile=-/opt/tradehub/.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1
ExecStart=/usr/bin/env node /opt/tradehub/build/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
