# tradehub/deploy/systemd/tradehub-parser.timer
[Unit]
Description=Периодический запуск парсера TradeHub (Telegram)

[Timer]
# Первый запуск через 2 минуты после загрузки, далее каждые 30 минут после завершения предыдущего прогона
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true
Unit=tradehub-parser.service

[Install]
WantedBy=timers.target
